# docker/images/php/supervisor/iiko-sync.conf
[program:iiko-sync]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --queue=iiko-sync --tries=1
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=8
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/supervisor/iiko-sync.log
stopwaitsecs=3600